# debian/postrm
#! /bin/sh
#
update-rc.d can2udp remove
service can2udp stop

