# debian/postinst
#! /bin/sh
#
systemctl enable can2udp
systemctl start can2udp


