# debian/postinst
#! /bin/sh
#
update-rc.d can2udp defaults
service can2udp restart


