# debian/postrm
#! /bin/sh
#
systemctl stop can2udp
systemctl disable can2udp

